# etc/systemd/system/cpupower.service
[Unit]
Description=CPU powersave
[Service]
Type=oneshot
ExecStart=/usr/local/bin/service_script
[Install]
WantedBy=multi-user.target